# 37-ipu6-psys.rules
ACTION=="add", SUBSYSTEM=="intel-ipu6-psys", MODE="0660", RUN+="/usr/bin/chgrp video '%E{DEVNAME}'"
